# supervisor-app.conf
[supervisord]
nodaemon=true

[program:cron]
command=/usr/sbin/cron -f
autostart=true
autorestart=true

[program:python-app]
command=python main.py
autostart=true
autorestart=true
